# coletor/com.unetvale.coletor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.unetvale.coletor</string>

  <!-- Caminho ABSOLUTO do venv + script. Ajuste se mover a pasta. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jhonicleyton/Desktop/Dashboard-Unetvale/dashboard-unetvale/coletor/.venv/bin/python</string>
    <string>/Users/jhonicleyton/Desktop/Dashboard-Unetvale/dashboard-unetvale/coletor/enviar.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jhonicleyton/Desktop/Dashboard-Unetvale/dashboard-unetvale/coletor</string>

  <!-- Horários: 08, 10, 12, 14, 16, 18 (se o Mac estava dormindo, roda ao acordar). -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/jhonicleyton/Desktop/Dashboard-Unetvale/dashboard-unetvale/coletor/coletor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jhonicleyton/Desktop/Dashboard-Unetvale/dashboard-unetvale/coletor/coletor.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
