# coletor/com.unetvale.coletor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Daemon do coletor (watcher). IMPORTANTE: rode de uma pasta FORA de
     ~/Desktop, ~/Documents, ~/Downloads — essas pastas são protegidas por TCC
     no macOS e o launchd NÃO consegue executar arquivos dentro delas (falha com
     exit 78 / EX_CONFIG, sem log). Local recomendado: ~/unetvale-coletor.

     Instalação:
       cp -R coletor/ ~/unetvale-coletor   (ou copie os .py + .env + requirements.txt)
       python3 -m venv ~/unetvale-coletor/.venv
       ~/unetvale-coletor/.venv/bin/pip install -r ~/unetvale-coletor/requirements.txt
       cp coletor/com.unetvale.coletor.plist ~/Library/LaunchAgents/   (ajuste os caminhos)
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.unetvale.coletor.plist

     O watcher roda continuamente (KeepAlive): cobre 08/10/12/14/16/18h mesmo
     após o Mac dormir e atende o botão "Atualizar" do app. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.unetvale.coletor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jhonicleyton/unetvale-coletor/.venv/bin/python</string>
    <string>/Users/jhonicleyton/unetvale-coletor/watcher.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/jhonicleyton/unetvale-coletor</string>
  <key>KeepAlive</key><true/>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/jhonicleyton/unetvale-coletor/coletor.log</string>
  <key>StandardErrorPath</key><string>/Users/jhonicleyton/unetvale-coletor/coletor.log</string>
</dict>
</plist>
